# templates/automatic-backup.service.j2
[Unit]
Description=Automatic backup {{ backup_name }}

[Service]
Type=simple
ExecStart=/bin/sh {{ backup_script_dir }}/automatic-backup-{{ backup_name }}.sh
User=aedu
Group=aedu
